# startup.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/testService.sh

[Install]
WantedBy=default.target
